# overlay/etc/s6-overlay/s6-rc.d/scheduled-tasks/finish
#!/command/with-contenv sh

set -euo pipefail
export BF_E="${PWD##*/}/$(basename ${0})"


#======================================================================================================================
# Do nothing - service will restart.
#======================================================================================================================

bf-svc-finish
