# overlay/etc/services.d/fetch-emails/finish
#!/usr/bin/with-contenv sh

set -euo pipefail


#======================================================================================================================
# Do nothing - service will restart.
#======================================================================================================================

bf-svc-finish -s ${PWD##*/}
