# overlay/etc/services.d/freescout-queue/finish
#!/usr/bin/with-contenv bash

set -euo pipefail


#======================================================================================================================
# Show helpful log message.
#======================================================================================================================

bf-svc-finish -s ${PWD##*/}
